# telegram_miniapp/deploy/trainer-miniapp-bot.service
[Unit]
Description=Trainer Telegram Mini App Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/trainer-miniapp/bot
EnvironmentFile=/etc/trainer-miniapp/bot.env
ExecStart=/usr/bin/python3 /opt/trainer-miniapp/bot/bot.py
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
