# backend/deploy/trainer-recommend-refresh.service
[Unit]
Description=Trainer: refresh coach recommendation if stale
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/trainer-miniapp/app
EnvironmentFile=/etc/trainer-miniapp/backend.env
ExecStart=/usr/bin/python3 /opt/trainer-miniapp/app/refresh_recommendation.py
User=root
