# backend/deploy/trainer-recommend-refresh.timer
[Unit]
Description=Trainer: morning refresh of the coach recommendation

[Timer]
# VPS runs Europe/Moscow — 06:30 local, so the card greets the athlete with a
# today-dated plan. Persistent=true catches up after VPS downtime.
OnCalendar=*-*-* 06:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
