# backend/deploy/trainer-db-backup.service
[Unit]
Description=Trainer: backup trainer.db (consistent snapshot + rotation)

[Service]
Type=oneshot
WorkingDirectory=/opt/trainer-miniapp/app
EnvironmentFile=/etc/trainer-miniapp/backend.env
ExecStart=/usr/bin/python3 /opt/trainer-miniapp/app/backup_db.py
User=root
